# ops-server/systemd/ops-dashboard.service
[Unit]
Description=Enablement Ops Dashboard
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=simple
User=ops
Group=ops
WorkingDirectory=/home/ops/enablement-framework/codespaces-framework/ops-server
ExecStart=/home/ops/ops-venv/bin/python -m dashboard.app
EnvironmentFile=/home/ops/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ops-dashboard

[Install]
WantedBy=multi-user.target
