# ops-server/systemd/ops-webhook.service
[Unit]
Description=Enablement Ops Webhook Server
After=network.target redis-server.service
Requires=redis-server.service

[Service]
Type=simple
User=ops
Group=ops
WorkingDirectory=/home/ops/enablement-framework/codespaces-framework/ops-server
ExecStart=/home/ops/ops-venv/bin/python -m uvicorn webhook.server:app --host 0.0.0.0 --port 8443 --log-level info
Restart=always
RestartSec=5
Environment=PYTHONPATH=/home/ops/enablement-framework/codespaces-framework/ops-server
EnvironmentFile=/home/ops/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ops-webhook

# Security
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/ops

[Install]
WantedBy=multi-user.target
